# systemd/duplicati.service
[Unit]
Description=Duplicati Backup Service

[Service]
Type=simple
ExecStart=/bin/duplicati
