# systemd/duplicati@.service
[Unit]
Description=Duplicati Backup Service for %i

[Service]
Type=simple
User=root
ExecStart=/bin/duplicati backup %i
